# ray_eap7.1/supervisord.conf
[supervisord]
nodaemon=true
[program:sshd]
command=/usr/sbin/sshd -D

[program:jboss-eap-7.1]
command=/home/jboss-eap-7.1/bin/domain.sh